# playbooks/roles/supernode/templates/bridge.service.j2
[Unit]
Description="Add br-nodes to batman"

[Service]
Type=oneshot
ExecStart=/bin/bash -c "batctl if add br-nodes"

[Install]
WantedBy=sys-subsystem-net-devices-br-nodes.device
